# deploy/supervisord.conf
[supervisord]
nodaemon=false

[program:celery]
environment=C_FORCE_ROOT="true"
command=celery -A backend_twous worker -l info
